# system/iwd.service
# vim: filetype=systemd
[Service]
ExecStartPre=/usr/bin/sleep 2
