# systemd/user/tmux-save.timer
[Unit]
Description=Save tmux session every ten minutes

[Timer]
OnCalendar=*:0/10

[Install]
WantedBy=timers.target
