# etc/systemd/system/iwd.service.d/override.conf
[Service]
; ExecStartPre=
; ExecStartPre=/usr/bin/sleep 3
; Restart=always

